# src/.config/systemd/user/task.service
[Unit]
Description=Sync taskwarrior

[Service]
ExecStart=/usr/bin/task sync

[Install]
WantedBy=multi-user.target
